# systemd/.config/systemd/user/newsboat.service
[Unit]
Description=Update RSS Feeds

[Service]
Type=oneshot
ExecStart=/usr/bin/newsboat -x reload
