# systemd/.config/systemd/user/newsboat.timer
[Unit]
Description=Update RSS Feeds

[Timer]
OnBootSec=2m
OnUnitActiveSec=6h
Unit=newsboat.service

[Install]
WantedBy=timers.target
